# deployment/systemd/buddhas-watch.service
[Unit]
Description=Buddhas-Watch CSI Defense Monitor
Documentation=https://github.com/ArcMike0430/Buddhas-Watch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=buddhas
Group=buddhas
WorkingDirectory=/opt/buddhas-watch/python/csi_monitor
ExecStart=/usr/bin/python3 /opt/buddhas-watch/python/csi_monitor/csi_defense.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=buddhas-watch

# Resource limits
LimitNOFILE=65536
Nice=-5

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
